# files/etc/uci-defaults/99-pisowifi-setup
#!/bin/sh

# Set up uhttpd to prefer index.php or handle redirection
# We will replace /www/index.html with a smart PHP redirector if it exists

if [ -f /www/index.html ]; then
    mv /www/index.html /www/index.html.bak
fi

cat << 'EOF' > /www/index.php
<?php
$portal = 'http://10.0.0.1/pisowifi/';
$o = [];
exec("uci -q get pisowifi.general.portal_url", $o, $ret);
if ($ret === 0 && isset($o[0]) && $o[0] !== '') {
    $portal = $o[0];
    if (substr($portal, -1) !== '/') {
        $portal .= '/';
    }
}
$host = $_SERVER['HTTP_HOST'];
// Adjust these IPs to match your router's LAN IP
$router_ips = ['openwrt.lan'];

// Check if the request is for the router admin interface
$is_admin = false;
foreach ($router_ips as $ip) {
    if (strpos($host, $ip) !== false) {
        $is_admin = true;
        break;
    }
}

if (!$is_admin) {
    // Captive portal user
    // Use absolute URL to force captive portal detection on some devices
    header("Location: " . $portal, true, 302);
    exit;
}

// Admin user - redirect to LuCI
header("Location: /cgi-bin/luci");
exit;
?>
EOF
 
# Create CNA endpoints
cat << 'EOF' > /www/generate_204.php
<?php
$portal = 'http://10.0.0.1/pisowifi/';
$o = [];
exec("uci -q get pisowifi.general.portal_url", $o, $ret);
if ($ret === 0 && isset($o[0]) && $o[0] !== '') {
    $portal = $o[0];
    if (substr($portal, -1) !== '/') {
        $portal .= '/';
    }
}
header("Location: " . $portal, true, 302);
exit;
?>
EOF

cat << 'EOF' > /www/hotspot-detect.php
<?php
$portal = 'http://10.0.0.1/pisowifi/';
$o = [];
exec("uci -q get pisowifi.general.portal_url", $o, $ret);
if ($ret === 0 && isset($o[0]) && $o[0] !== '') {
    $portal = $o[0];
    if (substr($portal, -1) !== '/') {
        $portal .= '/';
    }
}
header("Location: " . $portal, true, 302);
exit;
?>
EOF

cat << 'EOF' > /www/redirect.php
<?php
$portal = 'http://10.0.0.1/pisowifi/';
$o = [];
exec("uci -q get pisowifi.general.portal_url", $o, $ret);
if ($ret === 0 && isset($o[0]) && $o[0] !== '') {
    $portal = $o[0];
    if (substr($portal, -1) !== '/') {
        $portal .= '/';
    }
}
header("Location: " . $portal, true, 302);
exit;
?>
EOF

echo "Captive Portal" > /www/connecttest.txt

# Configure uhttpd to support PHP and index.php
uci set uhttpd.main.index_page='index.php'
uci set uhttpd.main.error_page='/index.php'
uci -q del_list uhttpd.main.interpreter='.php=/usr/bin/php-cgi'
uci add_list uhttpd.main.interpreter='.php=/usr/bin/php-cgi'
uci -q delete uhttpd.main.lua_prefix
uci set uhttpd.main.lua_prefix='/cgi-bin/luci'
uci set uhttpd.main.lua_handler='/usr/lib/lua/luci/sgi/uhttpd.lua'
uci set uhttpd.main.cgi_prefix='/cgi-bin'
uci commit uhttpd
/etc/init.d/uhttpd restart

# Force DNS to router IP for captive detection
uci -q del_list dhcp.@dnsmasq[0].address='/#/10.0.0.1'
uci add_list dhcp.@dnsmasq[0].address='/#/10.0.0.1'
uci commit dhcp
/etc/init.d/dnsmasq restart

# Ensure permissions
chmod +x /etc/init.d/pisowifi
chmod +x /usr/bin/pisowifi-backend.php

# Enable and start service
/etc/init.d/pisowifi enable
/etc/init.d/pisowifi start

exit 0
