# files/etc/uci-defaults/99-pisowifi-setup
#!/bin/sh

# Set up uhttpd to prefer index.php or handle redirection
# We will replace /www/index.html with a smart PHP redirector if it exists

if [ -f /www/index.html ]; then
    mv /www/index.html /www/index.html.bak
fi

cat << 'EOF' > /www/index.php
<?php
$host = $_SERVER['HTTP_HOST'];
// Adjust these IPs to match your router's LAN IP
$router_ips = ['192.168.1.1', '10.0.0.1', 'openwrt.lan'];

// Check if the request is for the router admin interface
$is_admin = false;
foreach ($router_ips as $ip) {
    if (strpos($host, $ip) !== false) {
        $is_admin = true;
        break;
    }
}

if (!$is_admin) {
    // Captive portal user
    header("Location: /pisowifi/");
    exit;
}

// Admin user - redirect to LuCI
header("Location: /cgi-bin/luci");
exit;
?>
EOF

# Ensure permissions
chmod +x /etc/init.d/pisowifi
chmod +x /usr/bin/pisowifi-backend.php

# Enable and start service
/etc/init.d/pisowifi enable
/etc/init.d/pisowifi start

exit 0
